# app01_Wordpress_mysql/deploy/systemd/cad-edu-healthcheck.timer
[Unit]
Description=Run cad-edu-healthcheck.service every 5 minutes (PLAN.md Phase 6)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
