# app01_Wordpress_mysql/deploy/systemd/cad-edu-healthcheck.service
[Unit]
Description=CAD Edu Platform - container health/resource check (PLAN.md Phase 6)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# Adjust to wherever this repo's app01_Wordpress_mysql directory actually
# lives on the deployment host.
WorkingDirectory=/opt/cad-edu-platform/app01_Wordpress_mysql
ExecStart=/usr/bin/env bash -c 'scripts/health-check.sh && scripts/intrusion-watch.sh'
